# push-to-talk.service
[Unit]
Description=Push-to-Talk Dictation Service
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/share/push-to-talk/venv/bin/python %h/.local/share/push-to-talk/push-to-talk.py
Restart=on-failure
RestartSec=5
StartLimitIntervalSec=120
StartLimitBurst=10

# DISPLAY is auto-detected by push-to-talk.py
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/%U/bus

[Install]
WantedBy=graphical-session.target
