# push-to-talk-update.service
[Unit]
Description=Update Push-to-Talk from Git

[Service]
Type=oneshot
ExecStart=%h/.local/share/push-to-talk/update.sh
