# push-to-talk-update.timer
[Unit]
Description=Check for Push-to-Talk updates daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
